# deployment/production/crontab
MAILTO=""
HOME=/usr/local/www/apache24/data/creative-city-berlin.de
PATH=/usr/local/www/apache24/data/creative-city-berlin.de/bin:/usr/local/sbin:/usr/local/bin:/usr/iports/bin:/usr/iports/sbin:/sbin:/bin:/usr/sbin:/usr/bin
SHELL=/usr/local/bin/bash
# every five minutes
*/5     *       *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/send_mails.sh                             && curl -fsS --retry 3 https://hc-ping.com/fbf9b477-e064-4db8-8640-35eac4b7ba31 > /dev/null
# hourly
15      *       *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/import_articles.sh                        && curl -fsS --retry 3 https://hc-ping.com/39a8b232-4b29-4d67-beb0-f3ea01e57d81 > /dev/null
20      *       *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/import_bulletins.sh                       && curl -fsS --retry 3 https://hc-ping.com/d99bd1fb-90c6-44c9-9576-1325e906ddc3 > /dev/null
# daily
00      01      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/cleanup.sh                                && curl -fsS --retry 3 https://hc-ping.com/22f2fa19-dea5-4f6f-9794-0db39d1b1b7c > /dev/null
30      01      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/daily_tasks.sh                            && curl -fsS --retry 3 https://hc-ping.com/7238c274-6ead-4136-8525-185498683513 > /dev/null
15      02      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/import_jobs_from_berlinstartupjobs.sh     && curl -fsS --retry 3 https://hc-ping.com/67cfddad-59fc-47e8-ac45-9757ad865f44 > /dev/null
20      02      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/import_jobs_from_creativeset.sh           && curl -fsS --retry 3 https://hc-ping.com/f2827e97-683e-4f99-b7ef-986b51c9c97b > /dev/null
25      02      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/import_jobs_from_kulturmanagement.sh      && curl -fsS --retry 3 https://hc-ping.com/570cd0b9-f16b-47b4-beca-36622107c562 > /dev/null
35      02      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/import_jobs_from_musicjob.sh              && curl -fsS --retry 3 https://hc-ping.com/f7aef8ec-35e8-4356-b0fd-b2c99b31dd51 > /dev/null
40      02      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/import_jobs_from_theaterjobs.sh           && curl -fsS --retry 3 https://hc-ping.com/9ef73899-7149-45b9-b6cf-803c4dde5ec9 > /dev/null
45      02      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/import_jobs_from_museumsbund.sh           && curl -fsS --retry 3 https://hc-ping.com/b5fd2dd9-32b5-4034-a78f-c896741a32f7 > /dev/null
50      02      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/import_jobs_from_gruenderszene.sh         && curl -fsS --retry 3 https://hc-ping.com/fe8f426d-53d3-4e91-bffe-c2a0a64d74a0 > /dev/null
10      03      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/import_jobs_from_dasauge.sh               && curl -fsS --retry 3 https://hc-ping.com/d883a096-d0ae-446d-b899-503168945027 > /dev/null
30      03      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/import_jobs_from_jobwrk.sh                && curl -fsS --retry 3 https://hc-ping.com/4270b34c-a6c1-432d-ace9-985ce547d92e > /dev/null
00      05      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/rebuild_index.sh                          && curl -fsS --retry 3 https://hc-ping.com/3982e561-e8ec-41e1-8a87-6cd0216f8009 > /dev/null
15      06      *       *       *       /usr/local/www/apache24/data/creative-city-berlin.de/commands/backup_postgres_db.sh                     && curl -fsS --retry 3 https://hc-ping.com/d137f0d5-1a8c-40ab-b69f-ea87a1cbc548 > /dev/null
# weekly
00      00      *       *       0       /usr/local/www/apache24/data/creative-city-berlin.de/commands/check_links.sh                            && curl -fsS --retry 3 https://hc-ping.com/9a5c9bb4-81af-4dfb-a2f2-9b98d0d961d8 > /dev/null
